# roles/zfs/templates/zfs-scrub.timer.j2
[Unit]
Description=Monthly ZFS scrub timer
Documentation=man:zpool-scrub(8)

[Timer]
OnCalendar={{ zfs_scrub_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
