# roles/zfs/templates/zfs-scrub.service.j2
[Unit]
Description=ZFS scrub on all pools
Documentation=man:zpool-scrub(8)
ConditionACPower=true
After=zfs-import.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/zpool scrub -w tank
