# roles/nginxplus/templates/lego-renew.timer.j2
[Unit]
Description=Run lego renew twice daily

[Timer]
OnCalendar={{ renew_calendar }}
Persistent=true

[Install]
WantedBy=timers.target

