# roles/solrcloud/templates/solr.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Apache Solr
After=network.target

[Service]
Type=forking
User={{ solr_user }}
Group={{ solr_group }}
Environment=SOLR_INCLUDE={{ solr_install_dir }}/solr/bin/solr.in.sh
ExecStart={{ solr_install_dir }}/solr/bin/solr start
ExecStop={{ solr_install_dir }}/solr/bin/solr stop
Restart=on-failure
LimitNOFILE=65000
LimitNPROC=65000

[Install]
WantedBy=multi-user.target

