# roles/solrcloud/templates/solr.service.j2
[Unit]
Description=Apache Solr
After=network.target

[Service]
Type=forking
User={{ solr_user }}
Group={{ solr_group }}
Environment=SOLR_INCLUDE=/etc/default/solr.in.sh
ExecStart={{ solr_bin }} start -cloud -z {{ solr_zookeeper_hosts_string }}/{{ solr_znode }}
ExecStop={{ solr_bin }} stop
Restart=on-failure
LimitNOFILE=65000
LimitNPROC=65000
TimeoutSec=180

[Install]
WantedBy=multi-user.target
